# etc/profile.d/canary.sh
#!/bin/sh
#echo "processing ~/usr/etc/profile.d/* for shell \"$0\" in stack ${BASH_SOURCE[*]}"
